# fsscripts/php.conf
description "PHP Server Daemon"
author "TrilioData, Inc"

start on startup
stop on shutdown
respawn

script
 sudo php -S 0.0.0.0:80 -t /home/ubuntu/tvault-filemanager
end script
